# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid
user=root
childlogdir=/var/log/supervisor

[program:thumbor]
command=thumbor --conf=/app/thumbor.conf -l %(ENV_THUMBOR_LOG_LEVEL)s --processes=%(ENV_THUMBOR_NUM_PROCESSES)s
stdout_logfile=/var/log/supervisor/thumbor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/thumbor-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autorestart=true
startretries=3
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
autorestart=true
autostart=true
startretries=3
priority=10
user=root
stopwaitsecs=10
killasgroup=true
stopasgroup=true